# config/supervisord.conf
[supervisord]
nodaemon=true

[program:clockwork]
command=bin/start_clockwork.sh

[program:worker]
command=bin/start_worker.sh
